# deploy/jlpt.service
[Unit]
Description=JLPT 單字王 FastAPI App
After=network.target

[Service]
Type=simple
User=APP_USER_PLACEHOLDER
Group=APP_USER_PLACEHOLDER
WorkingDirectory=APP_DIR_PLACEHOLDER
ExecStart=APP_DIR_PLACEHOLDER/.venv/bin/python main.py
Restart=on-failure
RestartSec=5

Environment="PRODUCTION=true"
Environment="PORT=8000"

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=APP_DIR_PLACEHOLDER/data

StandardOutput=journal
StandardError=journal
SyslogIdentifier=jlpt

[Install]
WantedBy=multi-user.target
